# oracle-deploy/systemd/profit-soul.service
[Unit]
Description=Profit Live Soul Master
After=network.target ollama.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/profit
ExecStart=/bin/bash live-soul-master.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
